# CryptoBot/deploy/oracle/callbuyer.service
[Unit]
Description=CallBuyer - Momentum Call Buying (15% allocation)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/home/botuser/CallBuyer
ExecStart=/home/botuser/CallBuyer/.venv/bin/python -u main.py

# Restart policy
Restart=always
RestartSec=30
StartLimitIntervalSec=600
StartLimitBurst=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=callbuyer

[Install]
WantedBy=multi-user.target
